# pantryos/pantryos-addon/rootfs/etc/s6-overlay/s6-rc.d/pantryos-app/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: PantryOS (Node Edition)
# Starts the Node.js application server
# ==============================================================================
set -euo pipefail

export NODE_ENV="production"
export APP_HOST="0.0.0.0"
export APP_PORT="80"
export APP_DATA_FILE="/data/pantryos/state.json"
export APP_PUBLIC_DIR="/opt/pantryos/public"
export APP_BASE_PATH="$(bashio::addon.ingress_entry)"
export APP_CULTURE="$(bashio::config 'culture')"
export APP_CURRENCY="$(bashio::config 'currency')"

if bashio::config.has_value 'timezone'; then
    export APP_TIMEZONE="$(bashio::config 'timezone')"
fi

if bashio::config.has_value 'log_level'; then
    export APP_LOG_LEVEL="$(bashio::config 'log_level')"
else
    export APP_LOG_LEVEL="info"
fi

bashio::log.info 'Starting PantryOS Node Edition service'
exec node /opt/pantryos/server/server.js
