# pantryos/pantryos/rootfs/etc/s6-overlay/s6-rc.d/init-pantryos/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: PantryOS (Node Edition)
# Prepares persistent storage for the refactored Node.js application
# ==============================================================================
set -euo pipefail

DATA_DIR="/data/pantryos"
STATE_FILE="${DATA_DIR}/state.json"
DEFAULT_STATE="/opt/pantryos/data/default-state.json"
EMPTY_STATE="/opt/pantryos/data/empty-state.json"

bashio::log.info 'Preparing persistent storage for PantryOS Node Edition'

if ! bashio::fs.directory_exists "${DATA_DIR}"; then
    bashio::log.notice 'Creating data directory'
    mkdir -p "${DATA_DIR}"
fi

if ! bashio::fs.file_exists "${STATE_FILE}"; then
    if bashio::config.false 'demo_data'; then
        bashio::log.notice 'Initializing with an empty dataset'
        cp "${EMPTY_STATE}" "${STATE_FILE}"
    else
        bashio::log.notice 'Seeding demo dataset'
        cp "${DEFAULT_STATE}" "${STATE_FILE}"
    fi
fi

if bashio::fs.file_exists "${STATE_FILE}"; then
    chmod 600 "${STATE_FILE}"
fi

touch "${DATA_DIR}/.keep"
chmod 600 "${DATA_DIR}/.keep"
