# pantryos/pantryos/rootfs/etc/services.d/pantryos/run
#!/usr/bin/with-contenv bash
set -euo pipefail

cd /opt/pantryos/server

# Porta interna esposta per Ingress
export PORT=8099

# Dove leggere/scrivere lo stato persistente
export PANTRYOS_DATA_DIR="/data/pantryos"

echo "[run] Starting PantryOS on port ${PORT}"
exec node server.js
